# .config/systemd/user/fetchmail@.timer
[Unit]
Description=Sync mailboxes with %i

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target

